# Bitcoin_Node/node_btc/supervisior.conf
[program:bitcoind]
command=/usr/local/bin/bitcoind -conf=/.bitcoin/bitcoin.conf